# pump-centinela/com.pump.analisis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pump.analisis</string>
  <!-- LUN a JUE 21:30 — corre al terminar el dia, a medida que van llegando
       los checks. Mati los revisa a la manana siguiente de a poco. El jueves
       es el cierre. El domingo (com.pump.centinela) solo PIDE.
       launchd NO acepta varios weekdays en un solo dict: hay que usar un
       array de dicts, uno por dia (1=lun, 2=mar, 3=mie, 4=jue). -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/matiassancari/agentkit-coach/venv/bin/python</string>
    <string>/Users/matiassancari/pump-centinela/centinela.py</string>
    <string>--analisis</string>
    <string>--send</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/matiassancari/pump-centinela</string>
  <key>StandardOutPath</key>
  <string>/Users/matiassancari/pump-centinela/analisis.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/matiassancari/pump-centinela/analisis.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
